# deploy/clep-update.service
[Unit]
Description=CLEP auto-update (git fetch + restart if changed)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=@APP_USER@
Group=@APP_GROUP@
WorkingDirectory=@APP_DIR@
ExecStart=@APP_DIR@/deploy/update.sh
