# deploy/clep.service
[Unit]
Description=CLEP — Calculator for Laser Experimentation Process
After=network-online.target avahi-daemon.service
Wants=network-online.target

[Service]
Type=simple
User=@APP_USER@
Group=@APP_GROUP@
WorkingDirectory=@APP_DIR@
ExecStart=@APP_DIR@/.venv/bin/python @APP_DIR@/app.py
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
